# deploy/systemd-user/kiosk-night.timer
[Unit]
Description=Stop kiosk every morning

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
